# modules/services/unix/logging/analysis_alert_action_server/files/alert_actioner/lib/dynamic_ip.service
[Unit]
Description=Dynamic IP
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/dhclient ens3:1

[Install]
WantedBy=multi-user.target